# scripts/com.postfinder.embed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.postfinder.embed</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/mathias/Projects/Codebase/postfinder_dev_app/postfinder-dev-search && uv run uvicorn server_local:app --host 0.0.0.0 --port 8000 --workers 1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mathias/Projects/Codebase/postfinder_dev_app/postfinder-dev-search</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mathias/Projects/Codebase/postfinder_dev_app/postfinder-dev-search/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mathias/Projects/Codebase/postfinder_dev_app/postfinder-dev-search/logs/server.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
